# docker/service/goto/run
#!/usr/bin/env bash
cd /site && hugo -b "$HUGO_BASEURL" && cd /docker
exec nginx -g 'daemon off;'
